# systemd/sunset-nag.service
[Unit]
Description=sunset-nag: progressive Night Light nag at sunset with snooze
After=graphical-session.target
PartOf=graphical-session.target

# NOTE: install.sh writes the real unit into ~/.config/systemd/user/ with the
# ExecStart path filled in to wherever this repo lives. This file is a template
# for reference.

[Service]
Type=simple
WorkingDirectory=%h/auto-grayscale
ExecStart=%h/auto-grayscale/sunset-nag daemon
Restart=always
RestartSec=5
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=%t/bus DISPLAY=:0

[Install]
WantedBy=default.target
